# schedule.cron
@hourly /usr/bin/jsub -N cat -quiet /data/project/nokib-bot/pwb/bin/python3 /data/project/nokib-bot/cat.py
@daily  /usr/bin/jstart -N daily -once -quiet /data/project/nokib-bot/pwb/bin/python3 /data/project/nokib-bot/daily.py
*/5 * * * * /usr/bin/jsub -N editWar -quiet /data/project/nokib-bot/pwb/bin/python3 /data/project/nokib-bot/editWar.py