# schedule.cron
0 * * * * /usr/bin/jsub -N cat  -once -quiet /data/project/goodarticlebot/pwb/bin/python3 /data/project/goodarticlebot/cat.py
0 */6 * * *  /usr/bin/jstart -N ga -once -quiet /data/project/goodarticlebot/pwb/bin/python3 /data/project/goodarticlebot/ga.py
*/5 * * * * /usr/bin/jsub -N editWar -quiet /data/project/goodarticlebot/pwb/bin/python3 /data/project/goodarticlebot/editWar.py